# squashfs-root/etc/init.d/watchdog
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=99

#USE_PROCD=1

start() {
        /usr/bin/watchdog
        return 0
}

stop() {
        return 0
}

