# squashfs-root/usr/lib/opkg/info/base-files.list
/usr/sbin/killpidfile
/usr/sbin/list_dir.lua
/usr/sbin/dhcp_apclient_get_ip.sh
/sbin/setssid_fact
/etc/init.d/smartvpn
/etc/openwrt_version
/etc/config/server_country_mapping
/etc/banner.failsafe
/usr/sbin/wanlinkprobe
/lib/preinit/99_01_prepare_eeprom
/etc/resolv.conf
/etc/board.d/01_leds
/usr/sbin/ban_admin
/usr/sbin/recovery_info_sync.sh
/usr/sbin/update_vas_conf.lua
/etc/init.d/done
/bin/config_generate
/usr/sbin/do_net_diagose.lua
/var
/lib/ramips.sh
/etc/init.d/wifiradio
/lib/preinit/00_preinit.conf
/usr/sbin/autowancheck
/etc/init.d/system
/etc/preinit
/usr/sbin/iwpriv_macfilter.lua
/bin/reboot_safe.sh
/usr/sbin/diskeh
/usr/sbin/guestwifi.sh
/etc/localtime
/usr/sbin/parentalctl.sh
/lib/upgrade/platform.sh
/usr/sbin/vas_security.lua
/etc/shells
/bin/flash_check.sh
/lib/config_post_ota/parentalctl_post_ota.sh
/lib/preinit/70_initramfs_test
/etc/device_info
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/usr/sbin/mivpn.sh
/etc/config/smartvpn
/lib/preinit/10_indicate_preinit
/etc/init.d/smp
/usr/libexec/login.sh
/usr/sbin/dhcp_apclient.sh
/usr/sbin/sysapi
/lib/config_post_ota/locale_info_ota.lua
/usr/sbin/supervisord
/usr/sbin/startscene.lua
/usr/sbin/ddnsd
/etc/config/server_mapping
/usr/sbin/gensmartdns.sh
/lib/config_post_ota/vasinfo_post_ota.sh
/usr/sbin/logupload.lua
/etc/init.d/bootcount
/etc/config/macfilter
/lib/preinit/32_check_for_boottype
/usr/sbin/cr_ip_conflict.lua
/etc/config/system
/usr/sbin/smp.sh
/usr/sbin/sync_guest_bssid.lua
/usr/sbin/config_collection.sh
/sbin/wps
/etc/fstab
/usr/lib/os-release
/bin/boardupgrade.sh
/etc/sysctl.conf
/sbin/sysupgrade
/usr/sbin/nslookupc
/lib/preinit/99_10_run_init
/usr/sbin/ubifs_chk
/usr/sbin/nslookupt
/etc/modules.d/92-mt7615e
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/usr/sbin/getstainfo.sh
/usr/sbin/crontab_record_camera.lua
/usr/sbin/kill_plugin_process.sh
/usr/lib/lua/json.lua
/etc/shadow
/etc/rc.common
/sbin/setmac
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/02_sysinfo
/usr/sbin/logcat
/etc/board.d/99-default_network
/usr/sbin/ap_mode.sh
/lib/preinit/10_indicate_failsafe
/usr/share/xiaoqiang/xiaoqiang-defaults.txt
/usr/sbin/wifiap_mode.sh
/etc/iwevent.d/01-feedpush.sh
/etc/os-release
/etc/board.d/02_network
/sbin/ethstatus
/usr/sbin/sysapi.firewall
/usr/sbin/jason.sh
/etc/init.d/watchdog
/usr/sbin/iperf_script.lua
/lib/preinit/40_mount_jffs2
/etc/init.d/cgroup_init
/lib/dnsmasq/lease.d/10-dhcp-mitvbox.sh
/etc/config/diag
/usr/sbin/setsmp.sh
/lib/preinit/80_mount_root
/lib/preinit/81_urandom_seed
/etc/rc.local
/usr/sbin/set_wifi_auto_ap_mode.lua
/lib/functions/uci-defaults.sh
/usr/sbin/macfilterctl
/usr/sbin/security_check.sh
/usr/sbin/ipset_redirect.sh
/usr/sbin/bootinfo
/usr/sbin/netdig.sh
/etc/opkg/keys/f130394897215624
/sbin/getmac
/usr/sbin/record_camera_by_ip.lua
/usr/sbin/wifishare.sh
/lib/preinit/41_mount_ext_data
/etc/ethers
/etc/init.d/gpio_switch
/lib/functions/network.sh
/etc/init.d/sysctl
/lib/preinit/90_mount_bind_etc
/sbin/hotplug-call
/usr/sbin/runt
/etc/crontabs/root
/usr/sbin/getneighbor.sh
/lib/preinit/04_handle_checksumming
/lib/upgrade/ubnt.sh
/usr/sbin/crontab_rom.sh
/etc/init.d/cron_fixup
/lib/upgrade/fwtool.sh
/etc/config/country_mapping
/etc/services
/usr/sbin/get_router_kv.lua
/etc/config/vas_user
/lib/preinit/10_check_for_mtd
/lib/functions/system.sh
/bin/flash.sh
/usr/sbin/vasinfo_fw.sh
/lib/upgrade/libupgrade.sh
/bin/board_detect
/usr/sbin/dhcp_do_opt43_act.sh
/lib/upgrade/common.sh
/etc/init.d/mtd_kf_check
/lib/config_post_ota/dnsfixd_post_ota.sh
/etc/init.d/urandom_seed
/etc/opkg/distfeeds.conf
/lib/preinit/31_restore_nvram
/etc/init.d/netconfig
/usr/sbin/sta_rssi_range.sh
/etc/sysupgrade.conf
/etc/hotplug.d/firmware/10-rt2x00-eeprom
/etc/config/wifishare
/usr/sbin/proxy_thirdparty.sh
/lib/preinit/07_set_preinit_iface_ramips
/lib/upgrade/keep.d/base-files
/usr/sbin/ecos_upgrade.lua
/etc/parentalctl/parentalctl.empty
/lib/preinit/31_check_for_boottype
/usr/sbin/wwdog
/lib/preinit/40_run_failsafe_hook
/etc/group
/usr/sbin/startscene_crontab.lua
/etc/config/vas
/etc/hotplug.d/iface/40-ntpd
/lib/preinit/30_failsafe_wait
/lib/upgrade/stage2
/etc/uci-defaults/13_fix_group_user
/sbin/firstboot
/usr/sbin/showmacs
/bin/login.sh
/usr/sbin/config_post_ota.sh
/etc/hotplug.d/net/00-sysctl
/lib/config_post_ota/ipv6fixd_post_ota.sh
/etc/hosts
/lib/functions.sh
/usr/sbin/smartvpn.sh
/etc/remove.conf
/etc/custom_hosts
/usr/sbin/mem_compaction.sh
/etc/config/network
/usr/sbin/upcrash
/lib/functions/leds.sh
/lib/config_post_ota/wireless_cfg_update.lua
/lib/xq-misc/phy_switch.sh
/etc/smartvpn/proxy.txt
/usr/sbin/openvpn_deamon.sh
/usr/sbin/vip_device_push_act.lua
/etc/config/parentalctl
/etc/init.d/timezone
/etc/board.d/03_gpio_switches
/usr/sbin/ntpsetclock
/usr/sbin/iptv.sh
/rom/note
/etc/inittab
/sbin/fixup-mac-address
/etc/openwrt_release
/lib/preinit/05_mount_skip
/sbin/hwversion
/usr/share/xiaoqiang/xiaoqiang_version
/lib/config_post_ota/wifishare_post_ota.sh
/sbin/wifi
/etc/iptv/route
/lib/wifi/ralink_common.sh
/etc/profile
/sbin/flash_led
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/usr/share/xiaoqiang/xiaoqiang-reserved.txt
/etc/TZ
/usr/sbin/guestwifi_mkip.lua
/etc/config/misc
/lib/preinit/90_nvram_fixup
/etc/init.d/umount
/lib/functions/service.sh
/usr/sbin/dropwatch
/lib/preinit/01_preinit_do_ramips.sh
/usr/sbin/wifiap_platform.sh
/lib/upgrade/keep.d/base-files-essential
/etc/sysctl.d/10-default.conf
/sbin/wifiradio
/usr/sbin/log_collection.sh
/etc/passwd
/lib/preinit/02_default_set_state
/usr/sbin/speed_test.lua
/usr/sbin/wanip_check.sh
/etc/hotplug.d/usb/10-motion
/etc/mtab
/etc/uci-defaults/09_fix-checksum
/lib/config_post_ota/ustack_post_ota.sh
/usr/sbin/dhcpc_do_opt43_act.sh
/etc/init.d/boot
/usr/sbin/phyhelper
/lib/functions/boot.sh
/usr/sbin/config_pre_ota.sh
/sbin/urandom_seed
/lib/functions/preinit.sh
/usr/sbin/syslog_upload.lua
/usr/sbin/checkupgrade.lua
/bin/ipcalc.sh
/usr/sbin/wifishare_daemon.sh
/lib/firewall/ipset_redirect.loader
/etc/diag.sh
/usr/sbin/killpid
/usr/sbin/gpio
/usr/sbin/wifiap_common.sh
/usr/sbin/nettb
