# squashfs-root/usr/lib/opkg/info/miqos.list
/usr/lib/lua/miqos/rule_by_noifb.lua
/usr/sbin/miqosc
/usr/lib/lua/miqos/common.lua
/usr/sbin/miqosd
/etc/init.d/miqos
/lib/config_post_ota/ifb_config_post_ota.sh
/usr/sbin/miqosd_init.sh
/usr/sbin/qos_monitor
/lib/config_post_ota/miqos_cfg_update.lua
/usr/lib/lua/miqos.lua
/usr/lib/lua/miqos/rule_by_host.lua
/lib/dnsmasq/lease.d/10-miio-mitv-qos-notify.sh
/usr/lib/lua/miqos/rule_by_prio.lua
/etc/config/miqos_default
/usr/lib/lua/miqos/command.lua
/usr/lib/lua/miqos/rule_by_wangzhe.lua
/etc/config/miqos
/usr/sbin/qos_monitor_hwqos
/usr/lib/lua/miqos/rule_by_service.lua
