# squashfs-root/usr/lib/opkg/info/nginx.list
/etc/nginx/nginx.conf
/etc/nginx/fastcgi_params
/etc/nginx/naxsi_core.rules
/etc/config/nginx
/etc/nginx/koi-win
/etc/nginx/nginx_https.conf
/etc/nginx/cert.crt
/etc/nginx/win-utf
/usr/sbin/nginx
/etc/nginx/koi-utf
/etc/nginx/cert.key
/etc/nginx/nginx_http.conf
/etc/init.d/nginx
/etc/nginx/htdocs/wifishare.html
/etc/nginx/mime.types
/lib/upgrade/keep.d/nginx
/etc/nginx/fastcgi-proxy-tcp.conf
/etc/nginx/miwifi-webinitrd.conf
/etc/nginx/force-redirect-https.conf
