# squashfs-root/usr/lib/opkg/info/taskdaemon.list
/etc/init.d/taskdaemon
/etc/taskdaemon.d/linkstatus/02-wan_dog.sh
/sbin/link_status_call
/usr/sbin/taskdaemon
